# source/services/convert_varrun.service
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=Convert /var/run into a symbolic link to /run .
DefaultDependencies=no
WantsMountsFor=/var
AfterMountsFor=/var
Before=devd devpubd local-fs.target

[Service]
Type=oneshot
ExecStart=./helper
# This is a one-shot service that can be re-triggered by starting it again.
Restart=no

[Install]
StoppedBy=shutdown.target
WantedBy=sysinit.target
