# source/services/utmps-wtmpd.socket
## **************************************************************************
## For copyright and licensing terms, see the file named COPYING.
## **************************************************************************

[Unit]
Description=utmps wtmp socket

[Socket]
Backlog=2
#MaxConnections=2
ListenStream=%t/utmps/.wtmpd-socket
Accept=on

[Install]
WantedBy=workstation.target
